# redirector/systemd/redirector.service
[Unit]
Description=redirector

[Service]
ExecStart=/usr/bin/redirector
User=nginx
Group=nginx
